# galaxy/keyd-openrc/usb-gadget.initd

#!/sbin/openrc-run
description="Keyd service for usb gadget setup"
command="/usr/bin/keyd-usb-gadget.sh"

depend() {
	before keyd 
}
