# galaxy/zramen-openrc/zramen.initd
#!/sbin/openrc-run

description="manage zram-backed swap devices"

start() {
	/usr/bin/zramen make
}

stop() {
       /usr/bin/zramen toss
}

