# system/dinit-rc/rc.local
#!/bin/sh
# rc.local for Artix -- enter your commands here
# it should be run when starting local.target service
